# infrastructure/config/systemd/configbuddy-backup-healthcheck.service
[Unit]
Description=ConfigBuddy Backup Health Check Service
Documentation=https://github.com/configbuddy/configbuddy
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
Group=root

# Environment variables
EnvironmentFile=-/etc/configbuddy/backup.env

Environment="BACKUP_DIR=/var/backups/configbuddy"
Environment="BACKUP_LOG_DIR=/var/log/configbuddy/backups"
Environment="BACKUP_MAX_AGE_HOURS=25"
Environment="BACKUP_MIN_COUNT=3"
Environment="DISK_USAGE_THRESHOLD=85"
Environment="BACKUP_NOTIFICATION_ENABLED=false"

# Working directory
WorkingDirectory=/opt/configbuddy/infrastructure/scripts

# Execute health check script
ExecStart=/opt/configbuddy/infrastructure/scripts/backup-health-check.sh

# Logging
StandardOutput=append:/var/log/configbuddy/backups/health-check.log
StandardError=append:/var/log/configbuddy/backups/health-check-error.log

# Restart policy
Restart=no

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
CPUQuota=25%
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
